# systemd-user/moltbot-config-sync.path
[Unit]
Description=Watch Moltbot config for changes

[Path]
PathChanged=%h/.clawdbot/clawdbot.json

[Install]
WantedBy=default.target
